# conf/supervisor.d/munged.ini
[program:munged]
autostart=true
command=/usr/local/bin/prefix-log.sh /sbin/munged -F --key-file /etc/munge/munge.key --num-threads 10 --socket=/var/run/munge/munge.socket.2
environment=LD_LIBRARY_PATH=/usr/local
user=munge
group=munge
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes = 0
priority=60
