# conf/supervisor.d/slurmdbd.ini
[program:slurmdbd]
autostart=true
command=/usr/local/bin/prefix-log.sh /usr/sbin/slurmdbd -D -v
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes = 0
priority=90

[program:slurmdbd_init__oneshot]
autostart=true
command=bash -c "sleep 10 && if [[ `sacctmgr list cluster Format=cluster%%30` != *%(ENV_CLUSTER_SLURM_NAME)s* ]]; then sacctmgr -i add cluster %(ENV_CLUSTER_SLURM_NAME)s && supervisorctl restart slurmdbd && supervisorctl restart slurmctld; else echo 'SlurmDBD cluster entry found'; fi"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes = 0
priority=91
