# conf/supervisor.d/slurmd.ini
[program:slurmd]
autostart=true
command=/usr/sbin/slurmd -D -f /etc/slurm/slurm.conf -v
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes = 0
priority=70
