# conf/supervisor.d/slurmctld.ini
[program:slurmctld]
autostart=true
command=/usr/local/bin/prefix-log.sh /usr/sbin/slurmctld -D -f /etc/slurm/slurm.conf -v
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes = 0
priority=90
